# add-ons/resources/addons.path
[Unit]
Description=Trigger for addons.service
[Path]
PathExists=/opt/addons/addons-start
Unit=addons.service
[Install]
WantedBy=multi-user.target